# systemd/resource-agents-deps.target
[Unit]
Description=resource-agents dependencies
