# ldirectord/systemd/ldirectord.service.in
[Unit]
Description=Monitor and administer real servers in a LVS cluster of load balanced virtual servers
Documentation=man:ldirectord(8)

[Service]
ExecStart=@sbindir@/ldirectord start
ExecStop=@sbindir@/ldirectord stop
ExecReload=@sbindir@/ldirectord reload
PIDFile=/var/run/ldirectord.ldirectord.pid
Type=forking
KillMode=none

[Install]
WantedBy=multi-user.target
