# ldirectord/systemd/ldirectord.service.in
[Unit]
Description=Monitor and administer real servers in a LVS cluster of load balanced virtual servers

[Service]
ExecStart=@sbindir@/ldirectord start
ExecStartPost=/usr/bin/touch /var/lock/subsys/ldirectord
ExecStop=@sbindir@/ldirectord stop
ExecStopPost=/usr/bin/rm -f /var/lock/subsys/ldirectord
ExecReload=@sbindir@/ldirectord reload
PIDFile=/var/run/ldirectord.ldirectord.pid
Type=forking
KillMode=none

[Install]
WantedBy=multi-user.target
